# home/.config/sv/swaybg/run
#!/bin/bash
wallpaper=$(cat "$HOME/.config/swaybg/wallpaper")
exec swaybg -i "$wallpaper" -m fill
