# home/.config/sv/swaybg/run
#!/bin/bash
wallpaper="$HOME/Pictures/Wallpapers/wallhaven-g7rprd-1.jpg"
exec swaybg -i "$wallpaper"
