# home/.config/sv/swaybg/run
#!/bin/bash
wallpaper="$HOME/Pictures/Wallpapers/wallhaven-l88ker.jpg"
exec swaybg -i "$wallpaper" -m fill
