# home/.config/sv/openrgb/run
#!/bin/bash
openrgb --startminimized --noautoconnect
