# home/.config/sv/openrgb/run
#!/bin/bash
exec openrgb --startminimized --noautoconnect
