# home/.config/sv/pipewire/run
#!/bin/bash
noisetorch() 
{
  sleep 2
  command noisetorch -i &
}

noisetorch &
exec pipewire

