# home/.config/sv/pipewire/run
#!/bin/bash
setup() {
    sleep 1 && ~/.local/bin/pipewire-setup.sh
}
setup &
exec pipewire

